# DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop terminal-notify.service
systemctl disable terminal-notify.service

exit 0